# supervisor.conf
[program:bingo]
directory=/home/ubuntu/bingo/
command = /usr/bin/npm start
user = ubuntu
autostart=true
autorestart=true
stdout_logfile = /home/ubuntu/bingo/supervisor.log
redirect_stderr = true